# test.service
[Unit]
Description=aaaxxx

[Service]
Type=simple
ExecStart=/opt/test/a

[Install]
WantedBy=multi-user.target